# docker/conf/supervisor/c9.ini
[program:c9]
command=/root/.c9/node/bin/node /usr/src/c9sdk/server.js -l 0.0.0.0 -p 444 -a admin:THE_PASSWORD
autorestart = true