# docker/conf/supervisor/flood.ini
[program:flood]
user=www
directory=/usr/www/rtorrent/flood/
command=/usr/bin/npm start
autorestart = true