# docker/conf/supervisor/flood.ini
[program:flood]
directory=/root/rtorrent/flood/
command=/usr/bin/npm start
autorestart = true