# docker/conf/supervisor/rtorrent.ini
[program:rtorrent]
user=www
environment=HOME="/home/www",PWD="/home/www",LOGNAME="www",USER="www",TERM=xterm
command=/usr/bin/rtorrent
autorestart = true