# docker/conf/supervisor/rtorrent.ini
[program:rtorrent]
environment=HOME="/root",PWD="/root",LOGNAME="root",USER="root",TERM=xterm
command=/usr/bin/rtorrent
autorestart = true