# docker/conf/nginx-php.ini
[program:php-fpm]
command=/usr/sbin/php-fpm -c /etc/php-fpm.conf
autorestart = true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autorestart = true