# docker/conf/rtorrent.ini
[program:rtorrent]
user=www
command=/usr/bin/rtorrent
autorestart = true