# docker/conf/rtorrent.ini
[program:rtorrent]
user=www
environment=TERM=xterm
command=/usr/bin/rtorrent
autorestart = true