# docker/conf/transmission.ini
[program:transmission]
user=www
command=/usr/bin/transmission-daemon -f --config-dir /usr/www/transmission/config/
autorestart = true

[program:socat]
command=/usr/bin/socat TCP4-LISTEN:443,reuseaddr,fork,su=nobody TCP4:127.0.0.1:9091
autorestart = true